# aports/forwarding/forwarding.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="Forwarding service"
description="Forwarding support"

command=/sbin/forwarding
command_background=true

depend() {
	need net
	after firewall
}
