# aports/uds_passthru/uds_passthru.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="Unix Domain Socket pass-through"
description="Unix Domain Socket pass-through support"

command=/sbin/uds_passthru
command_background=true

depend() {
	need net
	after firewall
}
