# aports/busybox/crond.initd
#!/sbin/openrc-run

name="busybox $SVCNAME"
command="/usr/sbin/$SVCNAME"
pidfile="/run/$SVCNAME.pid"
command_args="$CRON_OPTS -f"
command_background=true

depend() {
	need localmount
	need logger
}
