# aports/wpa_passthru/wpa_passthru.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="WPA Supplicant pass-through"
description="Control socket pass-through support for WPA Supplicant"

command=/sbin/wpa_passthru
command_background=true

depend() {
   need wpa_supplicant
}
